# supervisor/couchdb.conf
[program:couchdb]
autorestart=true
autostart=true
command=/bin/bash -c "exec /usr/local/var/lib/couchdb"
environment=HOME=/usr/local/var/lib/couchdb
redirect_stderr=true
user=couchdb
